# src/configs/systemd/openclaw.service
[Unit]
Description=OpenClaw Gateway — Self-hosted AI Assistant Platform
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=__OPENCLAW_USER__
Group=__OPENCLAW_USER__
WorkingDirectory=__OPENCLAW_HOME__
EnvironmentFile=-__OPENCLAW_HOME__/.openclaw/.env
Environment=NODE_OPTIONS=--max-old-space-size=1536
ExecStart=__OPENCLAW_BIN__ gateway run
Restart=always
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=__OPENCLAW_HOME__

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw

[Install]
WantedBy=multi-user.target
